# cdrom-mount@.service
[Unit]
Description=Mount CDROM Drive on %i
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/cdrom-mount.sh %i
